# src/sysinit/supervisor/lighttpd.conf
[program:lighttpd]
command=/usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf
priority=200

autostart=true
autorestart=unexpected

redirect_stderr=true
stdout_logfile=/var/log/lighttpd.log
stdout_logfile_maxbytes=10MB

