# src/sysinit/supervisor/php5-fpm.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm -F -O
priority=100

autostart=true
autorestart=unexpected

redirect_stderr=true
stdout_logfile=/var/log/php5.log
stdout_logfile_maxbytes=10MB

