# src/sysinit/supervisor/php7-fpm.ini
[program:php7-fpm]
command=/usr/sbin/php-fpm7 -F -O
priority=100

autostart=true
autorestart=unexpected

redirect_stderr=true
stdout_logfile=/var/log/php7.log
stdout_logfile_maxbytes=10MB

